# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; Unix socket for supervisorctl. Without this trio of sections, supervisord
; runs but supervisorctl has nothing to connect to and you can only restart
; processes by killing them and letting autorestart re-spawn them.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:vllm]
; bash -c lets us expand $VLLM_EXTRA_ARGS at runtime — supervisord itself
; doesn't expand shell vars in command lines. Use VLLM_EXTRA_ARGS for
; --quantization fp8, --tensor-parallel-size, etc.
command=/bin/bash -c "exec /opt/vllm-venv/bin/python -m vllm.entrypoints.openai.api_server --model %(ENV_MODEL_REPO)s --host %(ENV_VLLM_HOST)s --port %(ENV_VLLM_PORT)s --max-model-len %(ENV_MAX_MODEL_LEN)s --gpu-memory-utilization %(ENV_GPU_MEMORY_UTILIZATION)s --enable-prefix-caching $VLLM_EXTRA_ARGS"
autostart=true
autorestart=true
priority=10
startsecs=60
startretries=3
stdout_logfile=/var/log/supervisor/vllm.log
stderr_logfile=/var/log/supervisor/vllm-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=
    HF_HOME="%(ENV_HF_HOME)s",
    VLLM_EXTRA_ARGS="%(ENV_VLLM_EXTRA_ARGS)s"

[program:compactor]
; Runs in its OWN venv (/opt/compactor-venv), decoupled from vLLM as of
; V2.0 Phase 3 so chromadb/fastembed can't disturb vLLM's torch pins.
command=/opt/compactor-venv/bin/uvicorn main:app
    --host %(ENV_COMPACTOR_HOST)s
    --port %(ENV_COMPACTOR_PORT)s
    --log-level info
directory=/opt/compactor
autostart=true
autorestart=true
priority=15
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/compactor.log
stderr_logfile=/var/log/supervisor/compactor-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
; Only HF_HOME is listed explicitly (the tokenizer cache location must be
; right). Every other COMPACTOR_*/RAG/FASTEMBED var is inherited from the
; container environment (Dockerfile ENV defaults, overridable via .env or
; the RunPod template) — supervisord children inherit the parent env, so
; listing them here would be redundant AND risk a crash if one were unset.
environment=
    HF_HOME="%(ENV_HF_HOME)s",
    MODEL_REPO="%(ENV_MODEL_REPO)s",
    VLLM_URL="%(ENV_VLLM_URL)s"

[program:openwebui]
command=/app/venv/bin/open-webui serve --host 0.0.0.0 --port %(ENV_OPENWEBUI_PORT)s
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/openwebui.log
stderr_logfile=/var/log/supervisor/openwebui-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=
    HOME="/app",
    DATA_DIR="%(ENV_DATA_DIR)s",
    WEBUI_SECRET_KEY="%(ENV_WEBUI_SECRET_KEY)s",
    OPENAI_API_BASE_URL="%(ENV_OPENAI_API_BASE_URL)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    ENABLE_OLLAMA_API="%(ENV_ENABLE_OLLAMA_API)s",
    ENABLE_OPENAI_API="%(ENV_ENABLE_OPENAI_API)s",
    WEBUI_AUTH="%(ENV_WEBUI_AUTH)s"

[eventlistener:processes]
command=/bin/bash -c "echo READY && while read line; do echo OK; done"
events=PROCESS_STATE
