# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:llama-server]
command=/usr/local/bin/llama-server
    --model %(ENV_MODEL_PATH)s
    --host %(ENV_LLAMA_HOST)s
    --port %(ENV_LLAMA_PORT)s
    --ctx-size %(ENV_LLAMA_CTX_SIZE)s
    --n-gpu-layers %(ENV_LLAMA_N_GPU_LAYERS)s
    --parallel %(ENV_LLAMA_PARALLEL)s
    --cont-batching
    --flash-attn auto
    --mlock
    --context-shift
autostart=true
autorestart=true
priority=10
startsecs=30
startretries=3
stdout_logfile=/var/log/supervisor/llama-server.log
stderr_logfile=/var/log/supervisor/llama-server-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB

[program:openwebui]
command=/app/venv/bin/open-webui serve --host 0.0.0.0 --port %(ENV_OPENWEBUI_PORT)s
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/openwebui.log
stderr_logfile=/var/log/supervisor/openwebui-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=
    HOME="/app",
    DATA_DIR="%(ENV_DATA_DIR)s",
    WEBUI_SECRET_KEY="%(ENV_WEBUI_SECRET_KEY)s",
    OPENAI_API_BASE_URL="%(ENV_OPENAI_API_BASE_URL)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    ENABLE_OLLAMA_API="%(ENV_ENABLE_OLLAMA_API)s",
    ENABLE_OPENAI_API="%(ENV_ENABLE_OPENAI_API)s",
    WEBUI_AUTH="%(ENV_WEBUI_AUTH)s"

[eventlistener:processes]
command=/bin/bash -c "echo READY && while read line; do echo OK; done"
events=PROCESS_STATE
