# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

; Unix socket for supervisorctl. Without this trio of sections, supervisord
; runs but supervisorctl has nothing to connect to and you can only restart
; processes by killing them and letting autorestart re-spawn them.
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; ---------------------------------------------------------------------------
; Restart policy (V2.3 Theme 3 review)
; ---------------------------------------------------------------------------
; Long-running services below use autorestart=true + startretries=3 + a
; meaningful startsecs. The important property: a service that is GENUINELY
; BROKEN ON BOOT (e.g. vLLM OOMs because the model doesn't fit) crashes
; before startsecs elapses, which supervisord counts as a failed start;
; after startretries it transitions to FATAL and STAYS there — visible via
; `supervisorctl status`, not fast-restart-looping and hiding the root cause
; (the lesson from the V1.9.x cascade). Detection + recovery for a FATAL
; service is documented in OPERATIONS.md.
;   - startsecs : how long it must stay up to count as "started"
;   - startretries : failed starts before giving up → FATAL
;   - autorestart=true : restart if it dies AFTER a successful start
; One-shot/daemon sidecars (selftest, backup) set their own policy below.

[program:vllm]
; bash -c lets us expand $VLLM_EXTRA_ARGS at runtime — supervisord itself
; doesn't expand shell vars in command lines. Use VLLM_EXTRA_ARGS for
; --tensor-parallel-size, an offline-FP8 checkpoint, etc.
; startsecs=60: model load is slow; an OOM crash lands well inside 60s, so
; 3 retries → FATAL in ~3 min makes a too-big-model deploy fail visibly
; fast instead of looping. (On A40 use the 12B — see RUNPOD_DEPLOY.md.)
command=/bin/bash -c "exec /opt/vllm-venv/bin/python -m vllm.entrypoints.openai.api_server --model %(ENV_MODEL_REPO)s --host %(ENV_VLLM_HOST)s --port %(ENV_VLLM_PORT)s --max-model-len %(ENV_MAX_MODEL_LEN)s --gpu-memory-utilization %(ENV_GPU_MEMORY_UTILIZATION)s --enable-prefix-caching $VLLM_EXTRA_ARGS"
autostart=true
autorestart=true
priority=10
startsecs=60
startretries=3
stdout_logfile=/var/log/supervisor/vllm.log
stderr_logfile=/var/log/supervisor/vllm-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=
    HF_HOME="%(ENV_HF_HOME)s",
    VLLM_EXTRA_ARGS="%(ENV_VLLM_EXTRA_ARGS)s"

[program:compactor]
; Runs in its OWN venv (/opt/compactor-venv), decoupled from vLLM as of
; V2.0 Phase 3 so chromadb/fastembed can't disturb vLLM's torch pins.
command=/opt/compactor-venv/bin/uvicorn main:app
    --host %(ENV_COMPACTOR_HOST)s
    --port %(ENV_COMPACTOR_PORT)s
    --log-level info
directory=/opt/compactor
autostart=true
autorestart=true
priority=15
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/compactor.log
stderr_logfile=/var/log/supervisor/compactor-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
; Only HF_HOME is listed explicitly (the tokenizer cache location must be
; right). Every other COMPACTOR_*/RAG/FASTEMBED var is inherited from the
; container environment (Dockerfile ENV defaults, overridable via .env or
; the RunPod template) — supervisord children inherit the parent env, so
; listing them here would be redundant AND risk a crash if one were unset.
environment=
    HF_HOME="%(ENV_HF_HOME)s",
    MODEL_REPO="%(ENV_MODEL_REPO)s",
    VLLM_URL="%(ENV_VLLM_URL)s"

[program:openwebui]
command=/app/venv/bin/open-webui serve --host 0.0.0.0 --port %(ENV_OPENWEBUI_PORT)s
directory=/app
autostart=true
autorestart=true
priority=20
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/openwebui.log
stderr_logfile=/var/log/supervisor/openwebui-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
environment=
    HOME="/app",
    DATA_DIR="%(ENV_DATA_DIR)s",
    WEBUI_SECRET_KEY="%(ENV_WEBUI_SECRET_KEY)s",
    OPENAI_API_BASE_URL="%(ENV_OPENAI_API_BASE_URL)s",
    OPENAI_API_KEY="%(ENV_OPENAI_API_KEY)s",
    ENABLE_OLLAMA_API="%(ENV_ENABLE_OLLAMA_API)s",
    ENABLE_OPENAI_API="%(ENV_ENABLE_OPENAI_API)s",
    WEBUI_AUTH="%(ENV_WEBUI_AUTH)s"

# V3.2: Whisper speech-to-text service. Runs in its OWN venv (/opt/whisper-venv)
# on STT_PORT (9000). OpenWebUI's STT engine posts audio here; the compactor is
# NOT involved (audio -> text only — transcribed text then flows through the
# compactor like any typed message). CPU by default (server.py) so it never
# competes with vLLM for VRAM. Disable the whole service per-pod with
# STT_ENABLED=false. priority=25 puts it after openwebui(20), before selftest(30)
# so the boot self-test's STT check finds it already up.
[program:stt]
command=/opt/whisper-venv/bin/python /opt/stt/server.py
directory=/opt/stt
autostart=%(ENV_STT_ENABLED)s
autorestart=true
priority=25
startsecs=10
startretries=3
stdout_logfile=/var/log/supervisor/stt.log
stderr_logfile=/var/log/supervisor/stt-error.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
; All WHISPER_*/STT_* config is read from the inherited container environment
; (Dockerfile ENV defaults, overridable via .env / RunPod template). Not listed
; here so an unset var can't crash supervisord's config load.

# V2.1 Phase 6 Step 2: boot-time self-test.
# One-shot: runs once after vllm/compactor/openwebui are up, logs PASS/FAIL
# to /var/log/supervisor/selftest.log, then exits. Failure does NOT take
# down the pod — operator just sees it on the next observability check.
# Disable per-pod via COMPACTOR_SELFTEST_ON_BOOT=false.
# priority=30 puts it after vllm(10), compactor(15), openwebui(20).
[program:selftest]
command=/opt/compactor-venv/bin/python /opt/compactor/selftest.py --on-boot --wait-for-ready
directory=/opt/compactor
priority=30
autostart=%(ENV_COMPACTOR_SELFTEST_ON_BOOT)s
autorestart=false
startsecs=0
startretries=0
exitcodes=0,1
stdout_logfile=/var/log/supervisor/selftest.log
stderr_logfile=/var/log/supervisor/selftest-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
; STORAGE_ROOT intentionally NOT set here — memory.py defaults to
; /data/openwebui/compactor when COMPACTOR_STORAGE_ROOT is unset, and
; supervisord would crash on config load if we referenced an unset env.
environment=
    COMPACTOR_URL="http://127.0.0.1:%(ENV_COMPACTOR_PORT)s",
    VLLM_URL="http://%(ENV_VLLM_HOST)s:%(ENV_VLLM_PORT)s",
    MODEL_REPO="%(ENV_MODEL_REPO)s"

# V2.3 Theme 1: periodic data-durability backups.
# Long-running daemon: snapshots webui.db + the compactor memory store on a
# fixed interval, verifies each archive (PRAGMA integrity_check + JSON
# parse), prunes to retention, logs to /var/log/supervisor/backup.log. A
# failed cycle is logged and the daemon keeps running. Disable per-pod via
# COMPACTOR_BACKUP_ENABLED=false. priority=40 (after selftest=30).
[program:backup]
command=/opt/compactor-venv/bin/python /opt/compactor/backup.py --daemon
directory=/opt/compactor
priority=40
autostart=%(ENV_COMPACTOR_BACKUP_ENABLED)s
autorestart=true
startsecs=5
stdout_logfile=/var/log/supervisor/backup.log
stderr_logfile=/var/log/supervisor/backup-error.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
; All other backup config (dir, interval, retain, db path) is read from the
; environment by backup.py with sane defaults — not referenced here, so an
; unset var can't crash supervisord's config load.
environment=
    DATA_DIR="%(ENV_DATA_DIR)s"

[eventlistener:processes]
command=/bin/bash -c "echo READY && while read line; do echo OK; done"
events=PROCESS_STATE
